# packages/server/systemd/webbot-reverse-tunnel.service
[Unit]
Description=WebBot Reverse SSH Tunnel
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
SyslogIdentifier=webbot-reverse-tunnel
StandardOutput=journal
StandardError=journal
ExecStart=%h/bin/webbot-tunnel.sh %h/.ssh/id_ed25519_webbot_cloud root@182.43.86.126 -R 127.0.0.1:19090:127.0.0.1:9090
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
